# .config/autostart/sxhkd.desktop
[Desktop Entry]
Type=Application
Name=SXHKD
Description=Hotkey daemon to prevent middle click paste
Exec=sxhkd 

